# auction-price-service-ypm.service
[Unit]
Description=Auction Analytics Price Service (YPM only)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=exec
User=wavey
Group=wavey
WorkingDirectory=/Users/wavey/yearn/auction-analytics
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/Users/wavey/yearn/auction-analytics

# Run YPM price service only
ExecStart=/Users/wavey/yearn/auction-analytics/scripts/run_price_service.sh --pricer ypm --network electro

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=5

# Resource limits
LimitNOFILE=512
MemoryHigh=512M
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auction-price-ypm

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/Users/wavey/yearn/auction-analytics/logs

[Install]
WantedBy=multi-user.target